# webapp/config/dj_backup.crontab
30 3 * * * root /srv/myproject/config/dj_backup_rotated.sh
# Replace myproject if needed. Don't remove the empty line at the end of this file. It is required to run the cron job.
